# files/etc/uci-defaults/001-network-interface-config
uci set network.lan.ipaddr="172.17.17.1"

uci set network.VZW=interface
uci set network.VZW.proto="quectel"
uci set network.VZW.device="/dev/mhi_QMI0"
uci set network.VZW.apn="vzwinternet"
uci set network.VZW.auth="none"
uci set network.VZW.pdptype="ipv4v6"

uci set network.TM=interface
uci set network.TM.proto="quectel"
uci set network.TM.device="/dev/mhi_QMI01"
uci set network.TM.apn="fast.t-mobile.com"
uci set network.TM.auth="none"
uci set network.TM.pdptype="ipv4v6"

uci commit network
service network restart

